# systemd/monocons.service
[Unit]
Description=Rebuild monocons icon theme

[Service]
Type=oneshot
ExecStart=%h/.local/share/monocons/build.sh