# systemd/monocons.path
[Unit]
Description=Watch GTK accent changes

[Path]
PathChanged=%h/.config/gtk-4.0/

[Install]
WantedBy=default.target